# build-docker-images/wazuh-manager/config/etc/services.d/ossec-logs/run
#!/usr/bin/with-contenv sh

# dumping wazuh-manager.log to standard output
exec tail -F /var/wazuh-manager/logs/wazuh-manager.log