# wazuh-opendistro/config/etc/services.d/api/run
#!/usr/bin/with-contenv sh
echo >&2 "starting API"

exec /bin/node /var/ossec/api/app.js
