# wazuh-odfe/config/etc/services.d/ossec-logs/run
#!/usr/bin/with-contenv sh

# dumping ossec.log to standard output
exec tail -f /var/ossec/logs/ossec.log
